# run-deps/etc/profile.d/redhawk-sdrroot.csh
if (! $?SDRROOT ) then
    setenv SDRROOT /opt/redhawk/sdr
endif

